# web-server/awscli-sync-to-bucket.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task
    xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task" version="1.2">
    <RegistrationInfo>
        <Date>2019-10-18T04:03:23.5126233</Date>
        <Author>EC2AMAZ-MSAEP8B\Administrator</Author>
        <URI>\awscli-sync-to-bucket</URI>
    </RegistrationInfo>
    <Triggers>
        <CalendarTrigger>
            <Repetition>
                <Interval>PT5M</Interval>
                <Duration>P1D</Duration>
                <StopAtDurationEnd>false</StopAtDurationEnd>
            </Repetition>
            <StartBoundary>2019-10-18T00:00:00</StartBoundary>
            <ExecutionTimeLimit>P1D</ExecutionTimeLimit>
            <Enabled>true</Enabled>
            <ScheduleByDay>
                <DaysInterval>1</DaysInterval>
            </ScheduleByDay>
        </CalendarTrigger>
    </Triggers>
    <Principals>
        <Principal id="Author">
            <UserId>S-1-5-21-2373002294-500805120-4049372383-500</UserId>
            <LogonType>S4U</LogonType>
            <RunLevel>LeastPrivilege</RunLevel>
        </Principal>
    </Principals>
    <Settings>
        <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
        <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
        <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
        <AllowHardTerminate>true</AllowHardTerminate>
        <StartWhenAvailable>true</StartWhenAvailable>
        <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
        <IdleSettings>
            <StopOnIdleEnd>true</StopOnIdleEnd>
            <RestartOnIdle>false</RestartOnIdle>
        </IdleSettings>
        <AllowStartOnDemand>true</AllowStartOnDemand>
        <Enabled>true</Enabled>
        <Hidden>false</Hidden>
        <RunOnlyIfIdle>false</RunOnlyIfIdle>
        <WakeToRun>false</WakeToRun>
        <ExecutionTimeLimit>P1D</ExecutionTimeLimit>
        <Priority>7</Priority>
    </Settings>
    <Actions Context="Author">
        <Exec>
            <Command>powershell.exe</Command>
            <Arguments>c:\revize\scripts\awscli-sync-to-bucket.ps1 -ExecutionPolicy Bypass</Arguments>
        </Exec>
    </Actions>
</Task>
